# challenges/web-crazy-casino/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:redis]
command=redis-server --daemonize no --bind 127.0.0.1 --port 6379
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log
priority=1

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
priority=2

[program:nodejs]
command=node /app/src/app.js
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/nodejs.log
stderr_logfile=/var/log/supervisor/nodejs_err.log
environment=NODE_ENV="production",JWT_SECRET="8pRUZUwkV6gGV767",REDIS_URL="redis://127.0.0.1:6379"
priority=3
