# challenges/sqli-lunionfaitlaforce/src/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld
priority=10
autorestart=true

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
priority=20
autorestart=true
